# os/live-build/config/includes.chroot/etc/systemd/system/phoenix-storage-policy.service
[Unit]
Description=Enforce Phoenix storage policy
After=local-fs.target phoenix-data.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/phoenix-storage-policy.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
