# os/live-build/config/includes.chroot/etc/systemd/system/phoenix-ui.service
[Unit]
Description=Phoenix-15 UI server
After=network.target phoenix-backend.service phoenix-data.service phoenix-sync.service

[Service]
Type=simple
EnvironmentFile=/etc/phoenix/phoenix.env
WorkingDirectory=/opt/phoenix/app
ExecStart=/usr/local/bin/phoenix-ui.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
